# icmpcfg
#!/bin/sh /etc/rc.common

START=99

boot() {
	# enable ping
	sysctl -w net.ipv4.ping_group_range="0 1"
}
